# etc/s6-overlay/s6-rc.d/tailscaled/run
#!/command/with-contenv sh
exec /usr/local/bin/tailscaled --state=mem: --tun=userspace-networking
