# etc/s6-overlay/s6-rc.d/tailscale/up
/command/with-contenv /bin/sh -c "tailscale up --auth-key $TS_AUTHKEY --ssh"
